# monitoring.service
[Unit]
Description=test monitoring 
 
[Service]
ExecStart=/root/monitoring.sh
